# packages/prometheus-nginx/carbonio-prometheus-nginx-exporter.service
[Unit]
Description=Carbonio Prometheus NGINX exporter
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
Environment=SSL_VERIFY=false
EnvironmentFile=/etc/default/carbonio-prometheus-nginx-exporter
ExecStart=/usr/bin/carbonio-prometheus-nginx-exporter $ARGS
Restart=always

[Install]
WantedBy=multi-user.target
