# packages/prometheus-postgres/carbonio-prometheus-postgres-exporter.service
[Unit]
Description=Carbonio Prometheus Postgres exporter
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
Environment=SSL_VERIFY=false
EnvironmentFile=/etc/default/carbonio-prometheus-postgres-exporter
ExecStart=/usr/bin/carbonio-prometheus-postgres-exporter $ARGS
Restart=always

[Install]
WantedBy=multi-user.target
