# packages/prometheus-process/carbonio-prometheus-process-exporter.service
[Unit]
Description=Carbonio Prometheus Process Exporter
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
EnvironmentFile=/etc/default/carbonio-prometheus-process-exporter
ExecStart=/usr/bin/carbonio/carbonio-prometheus-process-exporter/carbonio-prometheus-process-exporter $ARGS
Restart=always

[Install]
WantedBy=multi-user.target
