# packages/prometheus-blackbox/carbonio-prometheus-blackbox-exporter.service
[Unit]
Description=Carbonio Prometheus Blackbox exporter
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
EnvironmentFile=/etc/default/carbonio-prometheus-blackbox-exporter
ExecStart=/opt/zextras/common/bin/carbonio-prometheus-blackbox-exporter $ARGS
Restart=always

[Install]
WantedBy=multi-user.target
