# packages/prometheus-consul/carbonio-prometheus-consul-exporter.service
[Unit]
Description=Carbonio Prometheus Consul exporter
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
EnvironmentFile=/etc/default/carbonio-prometheus-consul-exporter
ExecStart=/usr/bin/carbonio-prometheus-consul-exporter $ARGS
Restart=always

[Install]
WantedBy=multi-user.target
