# packages/prometheus-carbonio/carbonio-prometheus-carbonio-exporter.service
[Unit]
Description=Carbonio Prometheus Carbonio exporter
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
Environment=SSL_VERIFY=false
EnvironmentFile=/etc/default/carbonio-prometheus-carbonio-exporter
ExecStart=/usr/bin/carbonio-mailbox-exporter
Restart=always

[Install]
WantedBy=multi-user.target
