# packages/prometheus-alertmanager/carbonio-prometheus-alertmanager.service
[Unit]
Description=Carbonio Prometheus Alertmanager Service
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
EnvironmentFile=/etc/default/carbonio-prometheus-alertmanager
ExecStart=/opt/zextras/common/bin/carbonio-prometheus-alertmanager $ARGS
Restart=always

[Install]
WantedBy=multi-user.target
