# packages/prometheus-node/carbonio-prometheus-node-exporter.service
[Unit]
Description=Carbonio Prometheus Node Exporter
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
EnvironmentFile=/etc/default/carbonio-prometheus-node-exporter
ExecStart=/usr/bin/carbonio-prometheus-node-exporter $ARGS
Restart=always

[Install]
WantedBy=multi-user.target
