# packages/prometheus-mysqld/carbonio-prometheus-mysqld-exporter.service
[Unit]
Description=Carbonio Prometheus MySQL exporter
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
EnvironmentFile=/etc/default/carbonio-prometheus-mysqld-exporter
#ExecStartPre=/usr/bin/set-carbonio-prometheus-mysqld-exporter.sh
ExecStart=/usr/bin/carbonio-prometheus-mysqld-exporter $ARGS
Restart=always

[Install]
WantedBy=multi-user.target
