# packages/prometheus-openldap/carbonio-prometheus-openldap-exporter.service
[Unit]
Description=Carbonio Prometheus OpenLDAP exporter
Wants=network-online.target
After=network.target

[Service]
User=carbonio-prometheus
Group=carbonio-prometheus
Type=simple
Environment=SSL_VERIFY=false
EnvironmentFile=/etc/default/carbonio-prometheus-openldap-exporter
ExecStart=/usr/bin/carbonio-prometheus-openldap-exporter $ARGS
Restart=always

[Install]
WantedBy=multi-user.target
